# rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: AdGuard Home
# Runs the Nginx daemon
# ==============================================================================

# Wait for adguard to become available
bashio::net.wait_for 8765 localhost 900

bashio::log.info "Starting NGinx..."

# Disable HA Authentication if front door is open

exec nginx
